# etc/services.d/mot_detect/run
#!/command/with-contenv bash
exec 2>&1
python3 /ext/mot.v.3.0.py $MOT_CAMERANAME $MONGO_USER $MONGO_PASS $MONGO_URL $ROOK_ACCESS_KEY $ROOK_SECRET $ROOK_BUCKET $ROOK_SUBFOLDER $ROOK_URL